# supervisord.conf
[supervisord]

[program:python]
command=bash -c 'python run.py'
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
